# dracut/30coreos-installer/coreos-installer.service
[Unit]
Description=CoreOS Installer
After=coreos-installer.target

[Service]
Type=simple
ExecStart=/usr/libexec/coreos-installer
#StandardInput=tty-force
StandardInput=tty
#TTYPath=/dev/console
TTYPath=/dev/tty2
TTYReset=yes
TTYVHangup=yes
#TTYVTDisallocate=yes

[Install]
WantedBy=default.target
