# systemd/coreos-installer-pre.target
[Unit]
Description=Prepare for CoreOS Installer Target

Requires=coreos-installer-disable-device-auto-activation.service
